# share/systemd/smartmet-webmon.service
[Unit]
Description=SmartMet Server browser dashboard (smwebmon)
Documentation=man:smwebmon(1) https://github.com/fmidev/smartmet-monitor
After=network.target

[Service]
Type=simple
EnvironmentFile=-/etc/sysconfig/smartmet-webmon
ExecStart=/usr/bin/smwebmon $SMWEBMON_OPTS
# smwebmon must run as the user that owns the smartmetd processes
# (typically `smartmet-server`, created by the smartmet-server RPM)
# so two things work:
#   * `kernel.perf_event_paranoid=2` lets the daemon profile the
#     smartmetd processes — perf record cross-uid is denied at any
#     paranoid level above 0, and we want to keep the kernel sysctl
#     at its RHEL default.
#   * `/var/log/smartmet/*-access-log` (owned by the daemon writer)
#     is readable without ACL trickery.
# Group= is intentionally omitted so systemd derives it from the
# user's primary group in /etc/passwd, which differs across sites
# (htj at FMI, smartmet-server elsewhere). Override via a drop-in
# (`sudo systemctl edit smartmet-webmon`) if your deployment uses
# a different operator account.
User=smartmet-server
Restart=on-failure
RestartSec=5

# Defensive cgroup limits. Headroom over the cost-analysis figures
# (~100 MB RSS, ~5 % of one core sustained) plus the perf samplers'
# burst load: with --perf default-on, perf record / offcputime-bpfcc
# / friends spawn periodically as descendants of smwebmon and share
# its quota. 200% of one core gives the recording windows breathing
# room without letting a runaway eat the host. Raise (or remove) in
# a drop-in if you have a hot box and want full elasticity:
#   sudo systemctl edit smartmet-webmon
MemoryMax=512M
CPUQuota=200%

# smwebmon must call perf_event_open(2) and use bcc-tools to
# profile smartmetd. Five kernel checks block this for an
# unprivileged daemon, each one fixed by a narrow capability that
# is strictly less broad than the CAP_SYS_ADMIN we already need —
# so granting them does not expand the worst-case attack surface
# beyond what installing this unit already implies:
#
#   * smartmetd has NoNewPrivileges=1 (its own systemd unit's
#     hardening), which makes it non-traceable except by callers
#     holding CAP_SYS_PTRACE — same-uid does NOT bypass this and
#     no value of kernel.perf_event_paranoid bypasses it. So
#     CAP_SYS_PTRACE is needed for "perf record -p smartmetd-pid".
#
#   * Raw tracepoints (sched:sched_wakeup, block:block_rq_issue)
#     and ftrace function tracepoints are denied to unprivileged
#     users at paranoid >= 0 — so wakeup, blockflame, and the
#     bcc-tools modes need CAP_SYS_ADMIN. (CAP_PERFMON would be
#     the narrower fit on kernel >= 5.8, but RHEL 8's 4.18 doesn't
#     expose it; CAP_SYS_ADMIN is the next step up.)
#
#   * The bcc-tools that create kprobes (offcputime-bpfcc,
#     biolat-bpfcc, runqlat-bpfcc) write to
#     /sys/kernel/debug/tracing/kprobe_events. /sys/kernel/debug is
#     mode 700 root:root by default on RHEL, and CAP_SYS_ADMIN does
#     NOT bypass DAC checks — so traversing into that directory
#     requires CAP_DAC_READ_SEARCH (narrower than CAP_DAC_OVERRIDE;
#     read+traverse only, not write).
#
#   * perf_event_open's mmap path mlocks the per-CPU ring buffer.
#     Without CAP_IPC_LOCK an unprivileged caller is constrained
#     by RLIMIT_MEMLOCK and kernel.perf_event_mlock_kb (default
#     516 KB — too small for off-CPU's stack ring), and the call
#     fails with "mmap: Operation not permitted".
#
#   * /proc/kallsyms returns zeros to non-root readers when
#     kernel.kptr_restrict >= 1 (the RHEL default is 1) and perf
#     bails out at exit 255 because it can't locate the kernel
#     reference relocation symbol. CAP_SYSLOG bypasses
#     kptr_restrict=1 (but not =2 — that one needs the sysctl
#     lowered host-wide). On hardened hosts running kptr_restrict=2
#     the operator additionally uncomments the kptr_restrict line
#     in /usr/lib/sysctl.d/99-smartmet-perf.conf.
#
# AmbientCapabilities only takes effect when NoNewPrivileges is
# unset, so we drop it from this unit deliberately. The remaining
# hardening (ProtectSystem, ProtectHome, PrivateTmp) still applies.
AmbientCapabilities=CAP_SYS_PTRACE CAP_SYS_ADMIN CAP_DAC_READ_SEARCH CAP_IPC_LOCK CAP_SYSLOG
CapabilityBoundingSet=CAP_SYS_PTRACE CAP_SYS_ADMIN CAP_DAC_READ_SEARCH CAP_IPC_LOCK CAP_SYSLOG
NoNewPrivileges=no
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes

# This unit is shipped DISABLED by default. The cost of running it
# 24x7 was discussed and found unjustified versus on-demand use.
# Operators run `sudo systemctl start smartmet-webmon` when they
# want it, then SSH-tunnel: `ssh -L 8765:localhost:8765 host`.
[Install]
WantedBy=multi-user.target
